# data/swaynav.service
[Unit]
Description=swaynav - keyboard-driven mouse control daemon
Documentation=https://github.com/zisonzishen0415-stack/swaynav
After=graphical.target
Requires=ydotoold.service

[Service]
Type=simple
ExecStart=%h/.local/bin/swaynav
Restart=on-failure

[Install]
WantedBy=graphical.target